# yocto/meta-scmservice/recipes-webserver/webserver/files/scmservice.service
[Unit]
Description=Example Service
After=network.target

[Service]
ExecStart=/bin/scmservice
Restart=always

[Install]
WantedBy=multi-user.target